# provisioning/systemd/benchd.service
[Unit]
Description=Benchmark Daemon
After=network.target
Wants=benchmark.slice

[Service]
Type=simple
ExecStart=/usr/local/bin/benchd
User=root
Group=root
Restart=on-failure
RestartSec=5

# Pin to housekeeping core — never steal benchmark CPU time
CPUAffinity=0

# Security hardening
ProtectSystem=strict
ReadWritePaths=/var/run /sys/fs/cgroup/benchmark.slice /sys/devices/system/cpu /etc/benchd
PrivateTmp=true
NoNewPrivileges=false

# Required capabilities
AmbientCapabilities=CAP_SYS_NICE CAP_SYS_ADMIN
CapabilityBoundingSet=CAP_SYS_NICE CAP_SYS_ADMIN CAP_DAC_OVERRIDE

# Socket setup
RuntimeDirectory=benchd
ExecStartPre=/bin/rm -f /var/run/benchd.sock

Environment=BENCHD_SOCKET=/var/run/benchd.sock
Environment=BENCHD_LOG_LEVEL=info

[Install]
WantedBy=multi-user.target
