# provisioning/systemd/benchd.service
[Unit]
Description=Benchmark Daemon
After=network.target
Wants=benchmark.slice

[Service]
Type=simple
ExecStart=/usr/local/bin/benchd
User=root
Group=root
Restart=on-failure
RestartSec=5

# Pin to housekeeping core — never steal benchmark CPU time
CPUAffinity=0

ProtectSystem=strict
ReadWritePaths=/run /var/run /sys/fs/cgroup/benchmark.slice /sys/devices/system/cpu /etc/benchd
PrivateTmp=true
NoNewPrivileges=false

AmbientCapabilities=CAP_SYS_NICE CAP_SYS_ADMIN CAP_CHOWN CAP_FOWNER
CapabilityBoundingSet=CAP_SYS_NICE CAP_SYS_ADMIN CAP_DAC_OVERRIDE CAP_CHOWN CAP_FOWNER

RuntimeDirectory=benchd
ExecStartPre=/bin/rm -f /run/benchd/benchd.sock

Environment=BENCHD_SOCKET=/run/benchd/benchd.sock
Environment=BENCHD_LOG_LEVEL=info

[Install]
WantedBy=multi-user.target
