# provisioning/iso/first-boot.service
[Unit]
Description=Benchmark Appliance First Boot Setup
After=multi-user.target
ConditionPathExists=!/var/lib/bench/.setup-complete

[Service]
Type=oneshot
ExecStart=/usr/local/bin/bench-setup

# After setup: disable console autologin. The installer already marks
# .setup-complete internally, so we only need to clean up autologin here.
ExecStartPost=/bin/bash -c 'rm -rf /etc/systemd/system/getty@.service.d /etc/systemd/system/getty@tty1.service.d'

StandardInput=tty
StandardOutput=tty
TTYPath=/dev/tty1
TTYReset=yes
TTYVHangup=yes
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
