# provisioning/ansible/roles/benchd/templates/benchmark.slice.j2
[Unit]
Description=Benchmark CPU Isolation Slice
Before=slices.target

[Slice]
AllowedCPUs={{ isolated_cores }}
AllowedMemoryNodes=0
Delegate=yes

[Install]
WantedBy=slices.target
